# package/trop_diary_api.service
[Unit]
Description=Trop Diary API - Logging Service

[Service]
Type=simple
ExecStart=/usr/bin/trop-diary-api start --port 8080
