# deb/postinst
#!/usr/bin/env bash

set -e
set -x

SERVICE="trop-diary-api.service"

systemctl enable "$SERVICE"
systemctl start "$SERVICE"
